# device/maix/maix-root/usr/lib/opkg/info/adbd_auth_service.list
/bin/adbd_auth_service
/etc/init.d/adbd
