# device/maix/maix-root/usr/lib/opkg/info/busybox-init-base-files.list
/usr/lib/libopencv_xobjdetect.so.4.1
/home/keygen/libdl-2.28.so
/etc/openwrt_version
/etc/banner.failsafe
/usr/lib/libopencv_calib3d.so
/usr/lib/libopencv_tracking.so.4.1
/usr/lib/libopencv_rgbd.so
/usr/lib/libopencv_highgui.so.4.1
/var
/etc/udev/rules.d/10-local.rules
/usr/lib/libopencv_face.so
/usr/lib/libopencv_tracking.so.4.1.0
/usr/bin/keygen
/home/res/logo.png
/usr/lib/libopencv_datasets.so.4.1
/usr/lib/libopencv_ml.so.4.1
/usr/lib/libopencv_gapi.so.4.1
/etc/init.d/S50telnet
/home/keygen/keygen.tmp
/usr/lib/libopencv_imgcodecs.so.4.1.0
/usr/lib/libopencv_bioinspired.so.4.1
/usr/lib/libopencv_dnn.so.4.1
/usr/lib/libopencv_structured_light.so
/usr/lib/libopencv_videoio.so.4.1.0
/etc/shells
/usr/lib/libopencv_videostab.so
/etc/device_info
/usr/lib/libopencv_fuzzy.so
/rdinit
/usr/lib/libopencv_ccalib.so
/usr/lib/libopencv_img_hash.so.4.1.0
/usr/lib/libopencv_plot.so.4.1
/usr/lib/libopencv_fuzzy.so.4.1.0
/usr/lib/libopencv_shape.so.4.1.0
/usr/lib/libopencv_bgsegm.so.4.1
/usr/lib/libopencv_video.so.4.1.0
/usr/lib/libopencv_highgui.so.4.1.0
/home/keygen/libgcc_s-2.28.so
/usr/bin/maixpy3_config.py
/usr/lib/libopencv_xfeatures2d.so.4.1
/usr/lib/libopencv_features2d.so.4.1
/etc/udhcpd.conf
/usr/lib/libopencv_hfs.so.4.1.0
/usr/lib/libopencv_quality.so.4.1
/usr/bin/monitor-selfboot.sh
/usr/lib/libopencv_stereo.so.4.1.0
/usr/lib/libopencv_aruco.so.4.1.0
/usr/lib/libopencv_datasets.so
/usr/lib/libopencv_ml.so
/usr/lib/libopencv_phase_unwrapping.so.4.1.0
/usr/lib/libopencv_superres.so.4.1
/usr/lib/libopencv_reg.so.4.1
/usr/lib/libopencv_gapi.so
/usr/lib/libopencv_shape.so.4.1
/usr/lib/libopencv_stereo.so.4.1
/etc/fstab
/usr/lib/libopencv_dpm.so.4.1.0
/usr/lib/libopencv_freetype.so
/usr/lib/libopencv_line_descriptor.so.4.1
/usr/lib/libopencv_optflow.so
/usr/lib/libopencv_structured_light.so.4.1
/usr/lib/libopencv_flann.so.4.1.0
/usr/lib/libopencv_ccalib.so.4.1.0
/etc/banner
/usr/lib/libopencv_img_hash.so.4.1
/usr/lib/libopencv_hfs.so.4.1
/etc/init.d/rc.final
/usr/lib/libopencv_videoio.so.4.1
/etc/shadow
/usr/lib/libopencv_stereo.so
/usr/lib/libopencv_bgsegm.so.4.1.0
/etc/rc.common
/usr/lib/libopencv_video.so
/usr/lib/libopencv_optflow.so.4.1
/usr/lib/libopencv_stitching.so.4.1.0
/usr/lib/libopencv_imgcodecs.so.4.1
/home/res/_tmp_maix_opencv.cpython-38.so
/usr/lib/libopencv_features2d.so.4.1.0
/usr/lib/libopencv_rgbd.so.4.1.0
/usr/lib/libopencv_dpm.so.4.1
/usr/lib/libopencv_text.so
/usr/lib/libopencv_core.so.4.1
/usr/lib/libopencv_photo.so.4.1
/usr/lib/libopencv_stitching.so
/usr/lib/libopencv_face.so.4.1.0
/usr/lib/libopencv_bioinspired.so
/etc/init.d/S00mpp
/usr/lib/libopencv_dnn.so
/usr/lib/libopencv_quality.so
/run
/home/res/baars.ttf
/usr/lib/libopencv_bioinspired.so.4.1.0
/usr/lib/libopencv_superres.so
/lib/functions/uci-defaults.sh
/usr/lib/libopencv_photo.so
/etc/init.d/S41netparam
/usr/lib/libopencv_shape.so
/usr/lib/libopencv_datasets.so.4.1.0
/usr/lib/libopencv_ml.so.4.1.0
/lib/functions/network.sh
/usr/lib/libopencv_freetype.so.4.1
/usr/lib/libopencv_ximgproc.so.4.1
/usr/lib/libopencv_hfs.so
/usr/lib/libopencv_line_descriptor.so.4.1.0
/etc/hostname
/usr/lib/libopencv_flann.so.4.1
/etc/init.d/rc.preboot
/usr/bin/test_vivo
/usr/lib/libopencv_objdetect.so
/lib/functions/system.sh
/usr/lib/libopencv_dpm.so
/usr/lib/libopencv_freetype.so.4.1.0
/usr/lib/libopencv_core.so
/usr/lib/libopencv_xobjdetect.so.4.1.0
/usr/lib/libopencv_saliency.so.4.1
/usr/lib/libopencv_ccalib.so.4.1
/usr/lib/libopencv_plot.so
/usr/lib/libopencv_rgbd.so.4.1
/usr/lib/libopencv_xphoto.so.4.1.0
/etc/mdev.conf
/usr/lib/libopencv_face.so.4.1
/usr/lib/libopencv_img_hash.so
/etc/init.d/S50usb
/pseudo_init
/usr/lib/libopencv_xphoto.so.4.1
/usr/lib/libopencv_quality.so.4.1.0
/usr/lib/libopencv_imgproc.so.4.1.0
/etc/group
/usr/lib/libopencv_superres.so.4.1.0
/usr/lib/libopencv_phase_unwrapping.so
/etc/init.d/rc.modules
/etc/wpa_supplicant.conf
/usr/lib/libopencv_ximgproc.so
/usr/lib/libopencv_xphoto.so
/usr/lib/libopencv_text.so.4.1
/usr/lib/libopencv_stitching.so.4.1
/etc/init.d/S11dev
/etc/init.d/S01logging
/usr/lib/libopencv_phase_unwrapping.so.4.1
/usr/lib/libopencv_highgui.so
/lib/functions.sh
/usr/lib/libopencv_imgproc.so.4.1
/usr/lib/libopencv_objdetect.so.4.1
/lib/functions/leds.sh
/etc/init.d/rcK
/home/startup.py
/etc/init.d/rcS
/usr/lib/libopencv_videoio.so
/rom/note
/usr/lib/libopencv_dnn_objdetect.so
/etc/inittab
/etc/openwrt_release
/usr/lib/libopencv_aruco.so.4.1
/usr/bin/hotplug.sh
/usr/lib/libopencv_core.so.4.1.0
/etc/init.d/S20urandom
/usr/lib/libopencv_xobjdetect.so
/etc/opkg/keys/af22f7a88858c8e9
/etc/profile
/usr/lib/libopencv_reg.so
/usr/lib/libopencv_fuzzy.so.4.1
/home/keygen/libpthread-2.28.so
/usr/lib/libopencv_objdetect.so.4.1.0
/usr/lib/libopencv_surface_matching.so
/usr/lib/libopencv_dnn_objdetect.so.4.1
/etc/init.d/S02app
/usr/lib/libopencv_aruco.so
/usr/lib/libopencv_videostab.so.4.1.0
/usr/lib/libopencv_text.so.4.1.0
/usr/lib/libopencv_line_descriptor.so
/usr/bin/display
/lib/functions/service.sh
/usr/lib/libopencv_photo.so.4.1.0
/home/keygen/libc-2.28.so
/usr/lib/libopencv_xfeatures2d.so.4.1.0
/usr/lib/libopencv_dnn_objdetect.so.4.1.0
/usr/lib/libopencv_bgsegm.so
/etc/passwd
/usr/lib/libopencv_ximgproc.so.4.1.0
/usr/lib/libopencv_optflow.so.4.1.0
/usr/lib/libopencv_structured_light.so.4.1.0
/usr/lib/libopencv_xfeatures2d.so
/etc/mtab
/usr/lib/libopencv_surface_matching.so.4.1
/usr/lib/libopencv_features2d.so
/usr/lib/libopencv_saliency.so
/usr/lib/libopencv_video.so.4.1
/home/keygen/ld-2.28.so
/usr/lib/libopencv_imgcodecs.so
/usr/lib/libopencv_calib3d.so.4.1.0
/usr/lib/libopencv_imgproc.so
/lib/functions/preinit.sh
/usr/lib/libopencv_surface_matching.so.4.1.0
/usr/lib/libopencv_saliency.so.4.1.0
/usr/lib/libopencv_plot.so.4.1.0
/usr/lib/libopencv_dnn.so.4.1.0
/usr/lib/libopencv_videostab.so.4.1
/etc/init.d/S10udev
/usr/lib/libopencv_flann.so
/usr/lib/libopencv_calib3d.so.4.1
/usr/lib/libopencv_tracking.so
/usr/lib/libopencv_reg.so.4.1.0
/usr/lib/libopencv_gapi.so.4.1.0
/etc/init.d/S40network
